# rmrc-jetson-server/systemd/main-server.service
[Unit]
Description=Main Server (ZeroMQ Bridge + Coordination)
After=network-online.target lidar-stream.service
Wants=network-online.target
Requires=lidar-stream.service

[Service]
Type=simple
User=sosrobo
Group=dialout
WorkingDirectory=/home/sosrobo/rmrc-jetson-server
Environment="HOME=/home/sosrobo"
# ROS2 workspace setup + Python path
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source ~/rmrc_ws/install/setup.bash && cd /home/sosrobo/rmrc-jetson-server && python3 src/main_server.py'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=main-server

# USB permissions (dialout group)
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target
