# rmrc-jetson-server/systemd/unilidar-driver.service
[Unit]
Description=Unitree LiDAR Driver (ROS2 Launch)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sosrobo
Group=dialout
WorkingDirectory=/home/sosrobo
Environment="HOME=/home/sosrobo"
# ROS2 workspace setup: Humble + custom workspace
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source ~/rmrc_ws/install/setup.bash && ros2 launch unitree_lidar_ros2 launch.py'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unilidar-driver

# USB permissions (dialout group)
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target
